# chef/cookbooks/tw-webserver/templates/default/gunicorn.conf.erb
description "gunicorn application server running The Wreck"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit unlimited

setuid <%= node['wreck']['web_user'] %>
setgid <%= node['wreck']['web_group'] %>

env PATH=<%= node['wreck']['virtualenv'] %>/bin
chdir <%= node['wreck']['app_root'] %>/src
exec gunicorn --workers 10 --bind unix:/tmp/thewreck.sock -m 007 wsgi
